# build/sources/test.conf
[program:task]
#管理进程的命名
command=python3 -u /var/test.py
#执行的命令
stderr_logfile= "错误日志路径"
#错误日志输出路径
stdout_logfile="正常日志路径"
#日志输出路径
#directory=test
#命令执行的工作空间
autostart=true
#自动启动
user=root
#指定用户
autorestart=true
#自动重启
